# recipes-support/notify-led/files/init_led_complete.service
[Unit]
Description=Start LEDs when initialising
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/init_complete.sh
